# app-misc/yalp/files/parsers-initd-1
#!/sbin/runscript

command="/usr/bin/yalp-parsers"
command_args="${YALP_OPTS}"
command_background="1"
pidfile="/var/run/yalp-parsers.pid"
name="YALP parsers daemon"

depend() {
	use net logger
}

